# salt/roots/states/django/files/supervisor-django.conf
{% from "map.jinja" import global with context %}

[program:uwsgi]
user = vagrant
command = uwsgi --emperor {{ global.config_prefix }}/uwsgi.d/django.ini
autostart=true
autorestart=true
stderr_logfile = /var/log/uwsgi/err.log
stdout_logfile = /var/log/uwsgi/out.log
stopsignal=QUIT
